# novofon-bot.service
[Unit]
Description=NovoFon Voice Bot
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/NovoFon_Bot
Environment="PATH=/root/NovoFon_Bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Загрузка переменных из .env файла
EnvironmentFile=/root/NovoFon_Bot/.env

ExecStart=/root/NovoFon_Bot/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 9000 --workers 4
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=novofon-bot

[Install]
WantedBy=multi-user.target

